# google-cloud-ops-agent-fluent-bit.service.j2
[Service]
Environment=GOOGLE_SERVICE_CREDENTIALS=/etc/google/auth/application_default_credentials.json